# scripts/devsnow_thp.service
[Unit]
Description=THP data
After=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/app/scripts/THP_DATA.py
Restart=always
 
[Install]
WantedBy=multi-user.target

